# src/ayaled2.service
[Unit]
Description=Aya Neo LED manager

[Service]
ExecStart=/usr/bin/ayaled2

[Install]
WantedBy=multi-user.target
